# modules/puppet/files/etc/cron.d/puppet.cron
# THIS FILE IS MAINTAINED BY PUPPET!
30 * * * * root /usr/bin/puppet apply --modulepath /etc/puppet/modules /etc/puppet/manifests/nodes.pp 2>&1 > /tmp/puppet.log
